# ml/com.axiom.ml.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.axiom.ml.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/subhankarshukla/Desktop/claude code/axiom/valuation_app/ml/autorun.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/subhankarshukla/Desktop/claude code/axiom/valuation_app</string>

    <!-- Run Mon-Fri at 5:00pm local time (after US market close) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/subhankarshukla/Desktop/claude code/axiom/valuation_app/monitor_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/subhankarshukla/Desktop/claude code/axiom/valuation_app/monitor_launchd_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
